# .config/systemd/user/devmon.service
[Unit]
Description=Devmon Udevil Daemon

[Service]
ExecStart=/usr/bin/devmon

[Install]
WantedBy=default.target
