# .config/systemd/user/xmonitor.service
[Unit]
Description=Set the correct monitor and resolution

[Service]
Type=oneshot
ExecStart=%h/bin/xmonitor

[Install]
WantedBy=xsession.target
