# dot_config/systemd/user/xsession.target
[Unit]
Description=X session services
After=default.target
