# config/dot_config/systemd/user/docker-prune.service
[Unit]
Description=Prune unused Docker images

[Service]
Type=simple
ExecStart=/usr/bin/docker system prune -af --filter "until=72h"

[Install]
WantedBy=default.target
