# config/dot_config/systemd/user/dunst.service.tmpl
[Unit]
Description=Dunst notification daemon
After=xresources.service

[Service]
Type=dbus
BusName=org.freedesktop.Notifications
ExecStart={{ .nixBinPath }}/dunst

[Install]
WantedBy=xsession.target
