# config/dot_config/systemd/user/gpg-agent.service.d/override.conf.tmpl
[Service]
ExecStart=
ExecStart={{ .nixBinPath }}/gpg-agent --supervised
ExecReload=
ExecReload={{ .nixBinPath }}/gpgconf --reload gpg-agent
