# config/dot_config/systemd/user/gpg-agent.service.d/override.conf.tmpl
[Service]
ExecStart=
ExecStart={{ .nixProfile }}/bin/gpg-agent --supervised
ExecReload=
ExecReload={{ .nixProfile }}/bin/gpgconf --reload gpg-agent
